# autostart/skylon-remote.service
[Unit]
Description=SkylonRemoteApp - Rust HTTP server for Kettler elliptical trainer
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
# IMPORTANT: Change 'gilles' to your username before manual installation
# Or use the install-service.sh script which does this automatically
# Note: %h is a systemd variable that expands to the home directory of the User
User=gilles
WorkingDirectory=%h/Documents/SkylonRemoteApp
ExecStart=%h/.cargo/bin/cargo run --release
Restart=on-failure
RestartSec=10

# Log Management
# By default, logs are sent to journald
StandardOutput=journal
StandardError=journal

# Uncomment to reduce log verbosity (only warnings and errors)
#SyslogLevel=warning

# Uncomment to completely disable logs (NOT RECOMMENDED - makes debugging impossible)
#StandardOutput=null
#StandardError=null

[Install]
WantedBy=multi-user.target
